# .config/systemd/user/ydotoold.service
[Unit]
Description=ydotoold

[Service]
ExecStart=/usr/bin/ydotoold
Restart=always

[Install]
WantedBy=default.target
